# packages/parsagon/parsagon-0.4.2.tar.gz/parsagon-0.4.2/src/parsagon/supervisor.conf
[program:daphne]
command=bash -c "/home/ubuntu/parsagon/parsagon-local-server/src/parsagon/server/daphne.sh"
directory=/home/ubuntu/parsagon/parsagon-local-server/src/parsagon/server/
autostart=false
stopasgroup=true
user=ubuntu
stdout_logfile=/var/log/supervisor/daphne-stdout.log
stderr_logfile=/var/log/supervisor/daphne-stderr.log
startsecs=10

[program:celery]
command=bash -c "/home/ubuntu/parsagon/parsagon-local-server/src/parsagon/server/celery.sh"
directory=/home/ubuntu/parsagon/parsagon-local-server/src/parsagon/server/
autostart=false
user=ubuntu
stdout_logfile=/var/log/supervisor/celery-stdout.log
stderr_logfile=/var/log/supervisor/celery-stderr.log
startsecs=10
stopsignal=QUIT
stopasgroup=true
; Set Celery priority higher than default (999)
; so, if broker is supervised, it will start first.
priority=1000

[program:page-cache]
command=bash -c "/home/ubuntu/parsagon/parsagon-local-server/src/parsagon/server/page_cache.sh"
directory=/home/ubuntu/parsagon/parsagon-local-server/src/parsagon/server/
autostart=false
stopsignal=INT
stopasgroup=true
user=ubuntu
stdout_logfile=/var/log/supervisor/page-cache-stdout.log
stderr_logfile=/var/log/supervisor/page-cache-stderr.log
startsecs=10
